# systemd/claudegrid.service
[Unit]
Description=ClaudeGrid - Claude Code Session Visualizer
Documentation=https://github.com/your-repo/claudegrid
After=network.target

[Service]
Type=simple
User=%USER%
Group=%GROUP%
WorkingDirectory=%INSTALL_DIR%
ExecStart=%NODE_PATH% %INSTALL_DIR%/bin/claudegrid.js
Restart=on-failure
RestartSec=5

# Environment
Environment=NODE_ENV=production
Environment=CLAUDEGRID_PORT=%PORT%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claudegrid

# Security hardening
NoNewPrivileges=true
ReadWritePaths=%INSTALL_DIR%
PrivateTmp=no
ProtectSystem=no
ProtectHome=no
